# addon/roles/teleport-package/templates/teleport.service.j2
[Unit]
Description=Teleport Service

[Service]
Type=simple
ExecStart=teleport start --config=/root/.config/node_config.yaml
Restart=on-failure

[install]
WantedBy=multi-user.target
